# modules/nag/files/nagnag.cron
# /etc/cron.d/nagnag
#
# This file is managed and overwritten by Puppet. If you need
# to make changes, see the Puppetmaster and its configuration.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 *  * * *  root  /usr/local/bin/nagnag.sh
